# deploy/systemd/kiosk-app.service
[Unit]
Description=Kiosk Next.js application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__KIOSK_USER__
Group=__KIOSK_GROUP__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=-__INSTALL_DIR__/.env
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
