# deploy/systemd/kiosk-power.timer
[Unit]
Description=Check kiosk power schedule every minute

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
